# Dockerfiles/app/crontab
0 23 * * * /code/qiita.sh  >> /var/log/cron.log 2>&1
